# playbooks/network/templates/polaris-mail-egress.service.j2
[Unit]
Description=Polaris SMTP egress pinning rules
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/polaris-mail-egress.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
